# supervisord.conf
[program:apache]
command = apache2 -D FOREGROUND
autorestart = true